# devops/statbus-upgrade.service
[Unit]
Description=StatBus Upgrade Daemon
After=docker.service
Wants=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/statbus
ExecStart=/home/%i/statbus/sb upgrade daemon
Restart=always
RestartSec=30
RestartForceExitStatus=42
SuccessExitStatus=42
Environment=HOME=/home/%i

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=statbus-upgrade-%i

# Security hardening
# Note: NoNewPrivileges not set — daemon needs sudo for rsync (DB backup)
ProtectSystem=strict
ReadWritePaths=/home/%i/statbus /home/%i/statbus-backups /home/%i/statbus-maintenance

[Install]
WantedBy=multi-user.target
