# ops/statbus-upgrade-liveness@.timer
[Unit]
Description=StatBus Upgrade Liveness Observer timer (5-min tick)
# Drives statbus-upgrade-liveness@%i.service every 5 min (plan piece #7). The
# observer's N=30 min trip threshold ÷ this 5-min cadence gives ~6 checks per
# stall window — enough resolution to trip promptly at N without polling tighter
# than necessary. User-level — installed by ./sb install into
# ~/.config/systemd/user/ and enabled with `systemctl --user enable --now`.

[Timer]
# First check shortly after boot (the box may come up mid-wedge), then every 5
# min while active. OnUnitActiveSec re-arms relative to the LAST activation of
# the .service, so a slow liveness-check won't compound.
OnBootSec=2min
OnUnitActiveSec=5min
# Persistent: if the machine was off across a scheduled tick, run once on the
# next boot so a wedge that started before a reboot is still caught.
Persistent=true
Unit=statbus-upgrade-liveness@%i.service

[Install]
WantedBy=timers.target
